# packages/codex-opencode-go-router/examples/opencode-go-responses-proxy.service
[Unit]
Description=OpenCode Go Responses API adapter for Codex
After=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.env
Environment=OPENCODE_GO_PROXY_HOST=127.0.0.1
Environment=OPENCODE_GO_PROXY_PORT=61973
Environment=OPENCODE_GO_PROXY_DUMP_DIR=%h/.codex/opencode-go-proxy/dumps
ExecStart=/usr/bin/env node %h/codex-flows-public/packages/codex-opencode-go-router/src/router.mjs
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
